# deploy/systemd/deston.service
# =============================================================================
# DESTON — systemd unit (alternative to docker compose, for bare-metal Node).
# Install:
#   sudo cp deploy/systemd/deston.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now deston
# =============================================================================
[Unit]
Description=DESTON Next.js app
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/deston
EnvironmentFile=/opt/deston/.env
ExecStartPre=/usr/bin/npx prisma migrate deploy
ExecStart=/usr/bin/node /opt/deston/.next/standalone/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/deston/.next /opt/deston/node_modules /opt/deston/public /opt/deston/prisma

[Install]
WantedBy=multi-user.target
